# supervisor/3-php.conf
[program:php]
command = php-fpm -F
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

